# setup/content-gate.service
[Unit]
Description=Content Gate Dashboard Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/clawdbot/dashboard
ExecStart=/usr/bin/node /home/clawdbot/dashboard/server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=INBOX_DIR=/home/clawdbot/dashboard-inbox
Environment=SECOND_BRAIN_DIR=/home/clawdbot/second-brain

[Install]
WantedBy=default.target
